# crontab
#===============================================#
# Min | Hour | dayMth | Month | dayWk | Command #
#===============================================#
   0     6       *        *       *     update
